# systemd/user/gnome-keyring-daemon.service.d/unlock.conf
[Service]
ExecStart=
ExecStart=%h/bin/gnome-keyring-start
Restart=on-failure
RestartSec=5
